# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:openclaw]
command=/opt/openclaw/bin/openclaw gateway run --bind auto
user=openclaw
directory=/opt/openclaw
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/openclaw/openclaw.log
stderr_logfile=/var/log/openclaw/openclaw-error.log
environment=HOME="/home/openclaw",NODE_ENV="production",NODE_COMPILE_CACHE="/var/tmp/openclaw-compile-cache",OPENCLAW_NO_RESPAWN="1"

[program:telegram]
command=/opt/openclaw/bin/openclaw channels telegram
user=openclaw
directory=/opt/openclaw
autostart=true
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/var/log/openclaw/telegram.log
stderr_logfile=/var/log/openclaw/telegram-error.log
environment=HOME="/home/openclaw",NODE_ENV="production",NODE_COMPILE_CACHE="/var/tmp/openclaw-compile-cache",OPENCLAW_NO_RESPAWN="1"

[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
stdout_logfile=/var/log/openclaw/sshd.log
stderr_logfile=/var/log/openclaw/sshd-error.log

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stdout_logfile=/var/log/openclaw/cron.log
stderr_logfile=/var/log/openclaw/cron-error.log

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=false
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/caddy/caddy.log
stderr_logfile=/var/log/caddy/caddy-error.log

[program:fail2ban]
command=/usr/bin/fail2ban-server -xf start
autostart=false
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/var/log/openclaw/fail2ban.log
stderr_logfile=/var/log/openclaw/fail2ban-error.log

[group:openclaw-services]
programs=openclaw,telegram,sshd,cron,caddy,fail2ban
